# season-3-system-administration/episode-12-backup-recovery/solution/configs/systemd/backup-health-check.timer
# Episode 12: Backup & Recovery - Health Check Timer
# Liisa Kask — "Check backup health every 6 hours. Catch failures early."

[Unit]
Description=Backup Health Check Timer (Every 6 hours)
Documentation=https://github.com/kernel-shadows/episode-12
Requires=backup-health-check.service

[Timer]
# Schedule: Every 6 hours (00:00, 06:00, 12:00, 18:00)
OnCalendar=*-*-* 00,06,12,18:00:00

# Persistent
Persistent=true

# Run 2 minutes after boot
OnBootSec=2min

# Randomize (health checks can be staggered)
RandomizedDelaySec=5min

# Accuracy
AccuracySec=1min

# Unit to activate
Unit=backup-health-check.service

[Install]
WantedBy=timers.target

