# season-3-system-administration/episode-12-backup-recovery/solution/configs/systemd/backup-offsite.service
# Episode 12: Backup & Recovery - Offsite Backup Service
# Liisa Kask — "Offsite = survival. Fire, flood, Krylov can't reach remote backup."

[Unit]
Description=Offsite Backup Service (Remote Sync via SSH)
Documentation=https://github.com/kernel-shadows/episode-12
After=network-online.target
Wants=network-online.target

# Run after local backups finish
After=backup-full.service backup-incremental.service

[Service]
Type=oneshot
User=root
Group=root

# Execute offsite sync via backup_manager.sh
ExecStart=/usr/local/bin/backup_manager.sh backup_offsite

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=backup-offsite

# Timeout (network transfer may take time)
TimeoutStartSec=7200

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/backup /var/log /root/.ssh
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true

# Resource limits (network-bound, less CPU)
CPUQuota=40%
MemoryLimit=1G
TasksMax=50
IOWeight=50

# Nice priority
Nice=15

# Retry on failure (network issues common)
Restart=on-failure
RestartSec=300
StartLimitBurst=3
StartLimitIntervalSec=3600

[Install]
WantedBy=multi-user.target

