# season-3-system-administration/episode-12-backup-recovery/solution/configs/systemd/backup-incremental.timer
# Episode 12: Backup & Recovery - Incremental Backup Timer
# Liisa Kask — "Daily incremental Monday-Saturday. Fast. Efficient. Persistent."

[Unit]
Description=Incremental Backup Timer (Daily at 02:00, except Sunday)
Documentation=https://github.com/kernel-shadows/episode-12
Requires=backup-incremental.service

[Timer]
# Schedule: Monday-Saturday at 02:00 (Sunday = full backup)
OnCalendar=Mon..Sat *-*-* 02:00:00

# Persistent if system was off
Persistent=true

# Run 5 minutes after boot if missed
OnBootSec=5min

# Randomize (avoid collision with other services)
RandomizedDelaySec=10min

# Accuracy
AccuracySec=1min

# Unit to activate
Unit=backup-incremental.service

[Install]
WantedBy=timers.target

