# season-3-system-administration/episode-12-backup-recovery/solution/configs/systemd/backup-offsite.timer
# Episode 12: Backup & Recovery - Offsite Backup Timer
# Liisa Kask — "Sync to remote after local backups. Daily at 03:00. 3-2-1 rule."

[Unit]
Description=Offsite Backup Timer (Daily at 03:00)
Documentation=https://github.com/kernel-shadows/episode-12
Requires=backup-offsite.service

[Timer]
# Schedule: Every day at 03:00 (after local backups at 02:00)
OnCalendar=*-*-* 03:00:00

# Persistent (critical for offsite!)
Persistent=true

# Run 10 minutes after boot if missed
OnBootSec=10min

# Randomize (network transfers can be staggered)
RandomizedDelaySec=15min

# Accuracy
AccuracySec=5min

# Unit to activate
Unit=backup-offsite.service

[Install]
WantedBy=timers.target

