# season-3-system-administration/episode-12-backup-recovery/solution/configs/systemd/backup-health-check.service
# Episode 12: Backup & Recovery - Health Check Service
# Liisa Kask — "Monitor backup health. Silent failure = disaster waiting to happen."

[Unit]
Description=Backup Health Check Service
Documentation=https://github.com/kernel-shadows/episode-12
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root

# Execute health check via backup_manager.sh
ExecStart=/usr/local/bin/backup_manager.sh check_backup_health

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=backup-health-check

# Timeout (health check is fast)
TimeoutStartSec=300

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/log
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
LockPersonality=true

# Resource limits (very light)
CPUQuota=10%
MemoryLimit=256M
TasksMax=10

# Nice priority
Nice=19

# Exit status handling
# 0 = OK, 1 = WARNING, 2 = CRITICAL
SuccessExitStatus=0 1

# Alert on critical failure
# (In production, integrate with monitoring system)

[Install]
WantedBy=multi-user.target

