# season-3-system-administration/episode-12-backup-recovery/solution/configs/systemd/backup-incremental.service
# Episode 12: Backup & Recovery - Incremental Backup Service
# Liisa Kask — "Daily incremental = fast, efficient. Monday-Saturday."

[Unit]
Description=Incremental Backup Service (Daily Changes Only)
Documentation=https://github.com/kernel-shadows/episode-12
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root

# Execute incremental backup via backup_manager.sh
ExecStart=/usr/local/bin/backup_manager.sh backup_incremental

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=backup-incremental

# Timeout (incremental is faster than full)
TimeoutStartSec=1800

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/backup /var/log /home/student/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true

# Resource limits (incremental uses less resources)
CPUQuota=30%
MemoryLimit=1G
TasksMax=50
IOWeight=100

# Nice priority
Nice=10

# Success/failure handling
Restart=no
SuccessExitStatus=0

[Install]
WantedBy=multi-user.target

