# season-3-system-administration/episode-12-backup-recovery/solution/configs/systemd/backup-full.timer
# Episode 12: Backup & Recovery - Full Backup Timer
# Liisa Kask — "Weekly full backup = foundation. Every Sunday 02:00. No exceptions."

[Unit]
Description=Full Backup Timer (Weekly on Sunday 02:00)
Documentation=https://github.com/kernel-shadows/episode-12
Requires=backup-full.service

[Timer]
# Schedule: Every Sunday at 02:00
OnCalendar=Sun *-*-* 02:00:00

# If system was off during scheduled time, run on next boot
# CRITICAL: Liisa's rule - missed backups MUST run!
Persistent=true

# Run 5 minutes after boot if schedule was missed
OnBootSec=5min

# Randomize start time by up to 15 minutes
# Prevents "thundering herd" if multiple services scheduled at 02:00
RandomizedDelaySec=15min

# Accuracy (1 minute is fine for backups)
AccuracySec=1min

# Unit to activate
Unit=backup-full.service

[Install]
WantedBy=timers.target

