# season-3-system-administration/episode-12-backup-recovery/solution/configs/systemd/backup-full.service
# Episode 12: Backup & Recovery - Full Backup Service
# Liisa Kask — "Full backup every Sunday at 02:00. Foundation of 3-2-1 rule."

[Unit]
Description=Full Backup Service (Weekly Complete Backup)
Documentation=https://github.com/kernel-shadows/episode-12
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root

# Execute full backup via backup_manager.sh
ExecStart=/usr/local/bin/backup_manager.sh backup_full

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=backup-full

# Timeout (large backups may take time)
TimeoutStartSec=3600

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/backup /var/log /home/student/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true

# Resource limits (prevent backup from consuming all resources)
CPUQuota=50%
MemoryLimit=2G
TasksMax=100
IOWeight=100

# Nice priority (lower priority than interactive tasks)
Nice=10

# Success/failure handling
Restart=no
SuccessExitStatus=0
StandardInputText=

[Install]
WantedBy=multi-user.target

