# season-3-system-administration/episode-12-backup-recovery/starter/systemd/backup-full.timer
# Episode 12: Backup & Recovery - Full Backup Timer (STARTER)
# Liisa Kask — "Weekly full backup = foundation. Configure schedule correctly."

[Unit]
# TODO: Add description (hint: Full Backup Timer - Weekly on Sunday 02:00)
Description=

# TODO: Add documentation URL
Documentation=

# TODO: Which service does this timer activate? (hint: backup-full.service)
Requires=

[Timer]
# Schedule: When should full backup run?
# Syntax: DayOfWeek YYYY-MM-DD HH:MM:SS
# Examples:
#   Sun *-*-* 02:00:00        = Every Sunday at 02:00
#   Mon..Fri *-*-* 09:00:00   = Monday-Friday at 09:00
#   *-*-01 00:00:00           = 1st of month at midnight
#
# TODO: Set schedule for Sunday at 02:00
OnCalendar=

# Persistent: Run missed schedules on boot?
# If server was off during scheduled time, should backup run when server boots?
# CRITICAL for backups!
#
# TODO: Enable persistent (hint: true)
Persistent=

# If schedule was missed, run N minutes after boot
# Example: Server off on Sunday 02:00, boots Monday 10:00
# With OnBootSec=5min, backup runs at Monday 10:05
#
# TODO: Run how many minutes after boot if missed? (hint: 5min)
OnBootSec=

# Randomize start time by up to N minutes
# Prevents "thundering herd" if multiple services scheduled at 02:00
# Example: RandomizedDelaySec=15min → backup runs between 02:00-02:15
#
# TODO: Randomize by how much? (hint: 15min)
RandomizedDelaySec=

# Accuracy: How precise should timing be?
# 1min = timer can trigger within ±1 minute of scheduled time
# Allows systemd to batch timers for power efficiency
#
# TODO: Set accuracy (hint: 1min is fine for backups)
AccuracySec=

# Which unit to activate when timer triggers?
# TODO: Specify service unit (hint: backup-full.service)
Unit=

[Install]
WantedBy=timers.target

