# season-3-system-administration/episode-12-backup-recovery/starter/systemd/backup-incremental.service
# Episode 12: Backup & Recovery - Incremental Backup Service (STARTER)
# Liisa Kask — "Daily incremental = fast, efficient. Configure like full, but different command."

[Unit]
# TODO: Add description (hint: Incremental Backup Service - Daily Changes Only)
Description=

# TODO: Documentation URL
Documentation=

# TODO: Start after? (hint: network-online.target)
After=

# TODO: Wants dependency
Wants=

[Service]
Type=oneshot

# TODO: User (hint: root)
User=

# TODO: Group (hint: root)
Group=

# TODO: Command to execute (hint: backup_incremental instead of backup_full)
ExecStart=

# Logging
# TODO: StandardOutput
StandardOutput=

# TODO: StandardError
StandardError=

# TODO: SyslogIdentifier (hint: backup-incremental)
SyslogIdentifier=

# Timeout (incremental is faster than full)
# TODO: Timeout (hint: 1800 seconds = 30 minutes)
TimeoutStartSec=

# Security hardening (same as full backup)
# TODO: Fill in security settings (copy from backup-full.service)
NoNewPrivileges=
PrivateTmp=
ProtectSystem=
ProtectHome=
ReadWritePaths=
ProtectKernelTunables=
ProtectKernelModules=
ProtectControlGroups=
RestrictAddressFamilies=
RestrictNamespaces=
LockPersonality=

# Resource limits (incremental uses less resources than full)
# TODO: CPU quota (hint: 30% instead of 50%)
CPUQuota=

# TODO: Memory limit (hint: 1G instead of 2G)
MemoryLimit=

# TODO: Tasks (hint: 50 instead of 100)
TasksMax=

# TODO: I/O weight
IOWeight=

# TODO: Nice level
Nice=

# Success/failure handling
Restart=no
SuccessExitStatus=0

[Install]
WantedBy=multi-user.target

