# season-3-system-administration/episode-12-backup-recovery/starter/systemd/backup-full.service
# Episode 12: Backup & Recovery - Full Backup Service (STARTER)
# Liisa Kask — "Fill in TODOs. Configure correctly. Test thoroughly."

[Unit]
# TODO: Add description (What does this service do?)
Description=

# TODO: Add documentation URL
Documentation=

# TODO: Specify when service should start (after network online?)
After=

# TODO: Add Wants= dependency (network-online.target recommended)
Wants=

[Service]
# Type of service (oneshot = runs once and exits)
Type=oneshot

# TODO: Which user should run backup? (hint: root for /backup access)
User=

# TODO: Which group?
Group=

# TODO: What command to execute? (hint: /usr/local/bin/backup_manager.sh backup_full)
ExecStart=

# Logging
# TODO: Where should logs go? (hint: journal for centralized logging)
StandardOutput=
StandardError=

# TODO: Syslog identifier (how logs are tagged in journalctl)
SyslogIdentifier=

# Timeout (large backups may take time)
# TODO: How long to wait before timeout? (hint: 3600 seconds = 1 hour)
TimeoutStartSec=

# Security hardening
# TODO: Enable NoNewPrivileges (prevents privilege escalation)
NoNewPrivileges=

# TODO: Enable PrivateTmp (isolated /tmp directory)
PrivateTmp=

# TODO: Protect system directories (strict = read-only)
ProtectSystem=

# TODO: Protect home directories (read-only access)
ProtectHome=

# TODO: Allow write access to specific paths (hint: /backup /var/log /home/student/data)
ReadWritePaths=

# TODO: Enable additional security features
ProtectKernelTunables=
ProtectKernelModules=
ProtectControlGroups=

# TODO: Restrict network address families (hint: AF_UNIX AF_INET AF_INET6)
RestrictAddressFamilies=

# TODO: Restrict namespaces
RestrictNamespaces=

# TODO: Lock personality
LockPersonality=

# Resource limits (prevent backup from consuming all resources)
# TODO: Limit CPU usage (hint: 50% = half of one CPU core)
CPUQuota=

# TODO: Limit memory (hint: 2G = 2 gigabytes)
MemoryLimit=

# TODO: Limit number of tasks (hint: 100 processes/threads)
TasksMax=

# TODO: Set I/O weight (hint: 100 = normal priority)
IOWeight=

# Nice priority (lower priority than interactive tasks)
# TODO: Set nice level (hint: 10 = lower priority, range -20 to 19)
Nice=

# Success/failure handling
Restart=no
SuccessExitStatus=0

[Install]
WantedBy=multi-user.target

