# season-3-system-administration/episode-12-backup-recovery/starter/systemd/backup-incremental.timer
# Episode 12: Backup & Recovery - Incremental Backup Timer (STARTER)
# Liisa Kask — "Daily incremental Monday-Saturday. Sunday = full backup."

[Unit]
# TODO: Description
Description=

# TODO: Documentation
Documentation=

# TODO: Requires which service?
Requires=

[Timer]
# Schedule: Monday-Saturday at 02:00 (Sunday = full backup day)
# Syntax: Mon..Sat *-*-* 02:00:00
#
# TODO: Set schedule for Monday-Saturday at 02:00
OnCalendar=

# TODO: Enable persistent (hint: true)
Persistent=

# TODO: OnBootSec (hint: 5min)
OnBootSec=

# TODO: RandomizedDelaySec (hint: 10min)
RandomizedDelaySec=

# TODO: AccuracySec (hint: 1min)
AccuracySec=

# TODO: Which unit to activate?
Unit=

[Install]
WantedBy=timers.target

