# season-3-system-administration/episode-11-disk-management-lvm/solution/configs/disk-monitor.service
[Unit]
Description=Disk Space Monitoring Service — Episode 11
Documentation=https://github.com/kernel-shadows/episode-11
After=local-fs.target
# Ensure all filesystems are mounted before checking
Requires=local-fs.target

[Service]
# Service type: oneshot (run once and exit)
Type=oneshot

# Script to execute
ExecStart=/usr/local/bin/disk-monitor.sh

# Standard output/error → journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=disk-monitor

# Working directory
WorkingDirectory=/tmp

# User/Group
User=root
Group=root

# ============================================================================
# Security Hardening (Best Practices)
# ============================================================================

# Prevent privilege escalation
NoNewPrivileges=true

# Private /tmp (isolated from other processes)
PrivateTmp=true

# Read-only /usr, /boot, /efi
ProtectSystem=full

# Protect /home (read-only)
# Use 'false' if script needs to write to /home
ProtectHome=false

# Protect kernel tunables
ProtectKernelTunables=true

# Protect kernel modules
ProtectKernelModules=true

# Protect kernel logs
ProtectKernelLogs=true

# Protect /proc (hide processes from other users)
ProtectProc=invisible

# Hide /proc/[pid] except for own processes
ProcSubset=pid

# ============================================================================
# Resource Limits (Optional)
# ============================================================================

# CPU quota (10% of one CPU core)
# CPUQuota=10%

# Memory limit (256MB should be enough for monitoring script)
# MemoryLimit=256M

# Max number of tasks (threads/processes)
# TasksMax=10

# ============================================================================
# Restart Policy (Optional for oneshot, usually for Type=simple)
# ============================================================================

# Restart on failure (uncomment if needed)
# Restart=on-failure
# RestartSec=5min

# ============================================================================
# Timeout (Optional)
# ============================================================================

# Max time for script execution (default: 90s)
TimeoutStartSec=60s

[Install]
WantedBy=multi-user.target

# ============================================================================
# Installation Instructions:
# ============================================================================
# 1. Copy script to /usr/local/bin/:
#    sudo cp scripts/disk-monitor.sh /usr/local/bin/
#    sudo chmod +x /usr/local/bin/disk-monitor.sh
#
# 2. Copy service file to /etc/systemd/system/:
#    sudo cp configs/disk-monitor.service /etc/systemd/system/
#
# 3. Reload systemd daemon:
#    sudo systemctl daemon-reload
#
# 4. Test service manually:
#    sudo systemctl start disk-monitor.service
#
# 5. Check logs:
#    sudo journalctl -u disk-monitor.service -n 50
#
# 6. Enable service to run at boot (if needed):
#    sudo systemctl enable disk-monitor.service
#
# ============================================================================
# Troubleshooting:
# ============================================================================
# Check service status:
#   sudo systemctl status disk-monitor.service
#
# View recent logs:
#   sudo journalctl -u disk-monitor.service -f
#
# View errors only:
#   sudo journalctl -u disk-monitor.service -p err
#
# Test service without timer:
#   sudo systemctl start disk-monitor.service
#
# Disable service:
#   sudo systemctl disable disk-monitor.service
#   sudo systemctl stop disk-monitor.service
#
# ============================================================================

