# season-3-system-administration/episode-11-disk-management-lvm/solution/configs/disk-monitor.timer
[Unit]
Description=Disk Space Monitoring Timer — Episode 11
Documentation=https://github.com/kernel-shadows/episode-11
# Run after time synchronization (if using OnCalendar)
After=time-sync.target

[Timer]
# ============================================================================
# Schedule: Run every 30 minutes
# ============================================================================

# OnCalendar syntax: https://www.freedesktop.org/software/systemd/man/systemd.time.html
# Examples:
#   *:0/30          → Every 30 minutes (at :00 and :30)
#   *:15,45         → At :15 and :45 past every hour
#   *-*-* 00:00:00  → Daily at midnight
#   Mon *-*-* 09:00:00 → Every Monday at 9 AM

OnCalendar=*:0/30

# ============================================================================
# Boot Delay: Run 5 minutes after system boot
# ============================================================================
# Ensures filesystems are fully mounted before first check
OnBootSec=5min

# ============================================================================
# Persistent: Catch up missed runs
# ============================================================================
# If system was off when timer should have triggered → run on next boot
# Useful for laptops or systems with downtime
Persistent=true

# ============================================================================
# Accuracy: How precise should the timer be?
# ============================================================================
# AccuracySec=1s     → Very precise (default)
# AccuracySec=1m     → Within 1 minute window (saves power)
# AccuracySec=1h     → Within 1 hour window (very relaxed)
#
# For disk monitoring, 1 minute accuracy is fine (saves CPU cycles)
AccuracySec=1m

# ============================================================================
# Randomized Delay (Optional)
# ============================================================================
# Adds random delay (0 to specified value) before running
# Useful for distributed systems to avoid thundering herd problem
# RandomizedDelaySec=5min

# ============================================================================
# Unit to Activate
# ============================================================================
# By default, activates disk-monitor.service (same name as timer but .service)
# Can be explicitly specified:
# Unit=disk-monitor.service

[Install]
# Start timer when system reaches multi-user.target (normal boot)
WantedBy=timers.target

# ============================================================================
# Installation Instructions:
# ============================================================================
# 1. Copy timer file to /etc/systemd/system/:
#    sudo cp configs/disk-monitor.timer /etc/systemd/system/
#
# 2. Reload systemd daemon:
#    sudo systemctl daemon-reload
#
# 3. Enable timer (starts on boot):
#    sudo systemctl enable disk-monitor.timer
#
# 4. Start timer NOW:
#    sudo systemctl start disk-monitor.timer
#
# 5. Check timer status:
#    sudo systemctl status disk-monitor.timer
#
# 6. List all timers and when they'll run next:
#    systemctl list-timers disk-monitor.timer
#
# ============================================================================
# Viewing Timer Information:
# ============================================================================
# List all active timers:
#   systemctl list-timers
#
# Show next execution times:
#   systemctl list-timers --all
#
# Show detailed timer status:
#   systemctl status disk-monitor.timer
#
# Show when timer last triggered:
#   journalctl -u disk-monitor.timer -n 10
#
# Show logs from triggered service:
#   journalctl -u disk-monitor.service -n 50
#
# ============================================================================
# OnCalendar Syntax Examples:
# ============================================================================
# Every 15 minutes:
#   OnCalendar=*:0/15
#
# Every hour:
#   OnCalendar=hourly
#   OnCalendar=*:00:00
#
# Every day at 3 AM:
#   OnCalendar=daily
#   OnCalendar=*-*-* 03:00:00
#
# Every Monday at 9 AM:
#   OnCalendar=Mon *-*-* 09:00:00
#
# Every first day of month at midnight:
#   OnCalendar=*-*-01 00:00:00
#
# Multiple schedules (run on any):
#   OnCalendar=*:0/15    # Every 15 minutes
#   OnCalendar=daily     # AND once per day
#
# ============================================================================
# Testing OnCalendar Syntax:
# ============================================================================
# Test if syntax is valid and see next 5 execution times:
#   systemd-analyze calendar "*:0/30"
#   systemd-analyze calendar "Mon *-*-* 09:00:00"
#
# Example output:
#   Normalized form: *-*-* *:00,30:00
#   Next elapse: Sat 2025-10-11 13:00:00 UTC
#   (in UTC):    Sat 2025-10-11 13:00:00 UTC
#   From now: 4min left
#
# ============================================================================
# Troubleshooting:
# ============================================================================
# Timer not starting?
#   1. Check if timer is enabled:
#      systemctl is-enabled disk-monitor.timer
#   2. Check timer status:
#      systemctl status disk-monitor.timer
#   3. Check systemd logs:
#      journalctl -xe
#
# Service not running when timer triggers?
#   1. Check service status:
#      systemctl status disk-monitor.service
#   2. Check service logs:
#      journalctl -u disk-monitor.service -n 50
#   3. Test service manually:
#      sudo systemctl start disk-monitor.service
#
# Stop timer:
#   sudo systemctl stop disk-monitor.timer
#   sudo systemctl disable disk-monitor.timer
#
# ============================================================================

