# season-6-embedded-iot/episode-21-raspberry-pi/solution/configs/systemd/surveillance.service
[Unit]
Description=Surveillance Camera Stream - Episode 21
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/kernel-shadows/season-6-embedded-iot/episode-21-raspberry-pi/solution
ExecStart=/home/ubuntu/kernel-shadows/season-6-embedded-iot/episode-21-raspberry-pi/solution/camera_stream.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
