# season-3-system-administration/episode-10-processes-systemd/starter/systemd/shadow-backup.service
[Unit]
# KERNEL SHADOWS: Episode 10 — Backup Service (Triggered by Timer)
# Saint Petersburg, Russia
#
# TODO 1: Добавь Description для службы
# Подсказка: Description=KERNEL SHADOWS Automated Backup
Description=???

# TODO 2: Добавь документацию
# Подсказка: Documentation=https://github.com/kernel-shadows/episode-10
Documentation=???

# TODO 3: Настрой зависимости (After/Wants)
# Подсказка: После старта сети (network-online.target)
After=???
Wants=???

[Service]
# ============================================================================
# Service Configuration
# ============================================================================

# TODO 4: Выбери тип сервиса
# Подсказка: Type=oneshot (запускается один раз и завершается)
# oneshot используется для задач, которые выполняются и завершаются
# (backup, maintenance tasks), в отличие от simple (долго работающие демоны)
Type=???

# TODO 5: Укажи пользователя для запуска
# Подсказка: User=root (backup требует доступа ко всем файлам)
User=???

# TODO 6: Укажи рабочую директорию
# Подсказка: WorkingDirectory=/var/operations
WorkingDirectory=???

# ============================================================================
# Backup Command
# ============================================================================

# TODO 7: Укажи команду для выполнения backup
# Подсказка: ExecStart=/usr/local/bin/shadow-backup.sh
ExecStart=???

# TODO 8: Добавь pre-checks (проверки перед запуском)
# Подсказка 1: ExecStartPre=/usr/bin/test -d /var/operations
# Подсказка 2: ExecStartPre=/usr/bin/test -x /usr/local/bin/shadow-backup.sh
# Проверяем: существует ли директория? исполняемый ли скрипт?
ExecStartPre=???
ExecStartPre=???

# TODO 9: Добавь post-backup logging
# Подсказка: ExecStartPost=/usr/bin/logger -t shadow-backup "Backup completed successfully"
# logger отправляет сообщение в system log (journalctl)
ExecStartPost=???

# ============================================================================
# Timeouts
# ============================================================================

# TODO 10: Установи timeout (30 минут = 1800 секунд)
# Подсказка: TimeoutStartSec=1800s
# Backup больших файлов может занять время
TimeoutStartSec=???

# ============================================================================
# Resource Limits
# ============================================================================

# TODO 11: Ограничь CPU (50% одного ядра)
# Подсказка: CPUQuota=50%
# Backup не должен перегружать сервер
CPUQuota=???

# TODO 12: Ограничь память (512MB)
# Подсказка: MemoryLimit=512M
MemoryLimit=???

# TODO 13: Установи приоритет I/O (низкий)
# Подсказка: IOSchedulingClass=2 (best-effort), IOSchedulingPriority=7 (lowest)
# Backup не должен замедлять production I/O
IOSchedulingClass=???
IOSchedulingPriority=???

# ============================================================================
# Security
# ============================================================================

# TODO 14: Включи приватный /tmp
# Подсказка: PrivateTmp=true
# Изолирует временные файлы backup от других процессов
PrivateTmp=???

# ============================================================================
# Logging
# ============================================================================

# TODO 15: Настрой вывод логов в journal
# Подсказка: StandardOutput=journal, StandardError=journal
StandardOutput=???
StandardError=???

# TODO 16: Установи идентификатор для syslog
# Подсказка: SyslogIdentifier=shadow-backup
# Позволяет фильтровать логи: journalctl -t shadow-backup
SyslogIdentifier=???

[Install]
# Примечание: [Install] секция НЕ нужна для oneshot сервиса, который запускается через timer
# Timer сам активирует этот сервис по расписанию

# ============================================================================
# СПРАВКА: Как использовать после заполнения TODO
#
# 1. Скопировать в systemd:
#    sudo cp shadow-backup.service /etc/systemd/system/
#
# 2. Reload systemd:
#    sudo systemctl daemon-reload
#
# 3. Ручной запуск (тест):
#    sudo systemctl start shadow-backup.service
#
# 4. Проверка статуса:
#    sudo systemctl status shadow-backup.service
#
# 5. Просмотр логов:
#    sudo journalctl -u shadow-backup.service
#
# 6. Автоматический запуск:
#    Настраивается через shadow-backup.timer (см. следующий файл)
#
# ============================================================================
# ВОПРОСЫ ДЛЯ ПОНИМАНИЯ:
#
# 1. Почему Type=oneshot, а не Type=simple?
#    → oneshot для задач, которые завершаются (backup)
#    → simple для долго работающих демонов (web server, database)
#
# 2. Зачем ExecStartPre проверки?
#    → Гарантируют что все необходимые файлы/директории существуют
#    → Сервис не запустится если что-то не так
#
# 3. Почему CPUQuota=50%, а не 100%?
#    → Backup не должен перегружать production сервер
#    → Оставляем ресурсы для других сервисов
#
# 4. Почему User=root?
#    → Backup должен читать ВСЕ файлы системы
#    → Для обычных сервисов используй непривилегированного пользователя!
#
# ============================================================================
# Борис: "Service = что делать. Timer = когда делать.
#         Разделение ответственности. SystemD архитектура."
# ============================================================================

