# season-3-system-administration/episode-10-processes-systemd/solution/systemd/shadow-backup.service
[Unit]
# KERNEL SHADOWS: Episode 10 — Backup Service (Triggered by Timer)
# Saint Petersburg, Russia
#
# Description: Creates encrypted backup of Shadow operations data
Description=KERNEL SHADOWS Automated Backup
Documentation=https://github.com/kernel-shadows/episode-10
After=network-online.target
Wants=network-online.target

[Service]
# ============================================================================
# Service Configuration
# ============================================================================
Type=oneshot
# oneshot: Run once and exit (perfect for backup tasks)
# Timer will trigger this service periodically

# User/Group
User=root
# ^ Needs root to read all operational files

# Working Directory
WorkingDirectory=/var/operations

# ============================================================================
# Backup Command
# ============================================================================
ExecStart=/usr/local/bin/shadow-backup.sh

# Pre-checks
ExecStartPre=/usr/bin/test -d /var/operations
ExecStartPre=/usr/bin/test -x /usr/local/bin/shadow-backup.sh

# Post-backup verification
ExecStartPost=/usr/bin/logger -t shadow-backup "Backup completed successfully"

# ============================================================================
# Timeouts
# ============================================================================
# Max runtime: 30 minutes (backup large datasets)
TimeoutStartSec=1800s

# ============================================================================
# Resource Limits
# ============================================================================
# CPU: Max 50% (backup can be CPU-intensive)
CPUQuota=50%

# Memory: Max 512MB
MemoryLimit=512M

# I/O Priority: Low (don't impact production services)
IOSchedulingClass=2
IOSchedulingPriority=7

# ============================================================================
# Security
# ============================================================================
# Private /tmp
PrivateTmp=true

# ============================================================================
# Logging
# ============================================================================
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shadow-backup

[Install]
# No [Install] section needed — triggered by timer, not enabled directly
# See shadow-backup.timer for scheduling

# ============================================================================
# Usage:
#
# 1. Manual run:
#    sudo systemctl start shadow-backup.service
#
# 2. Check status:
#    sudo systemctl status shadow-backup.service
#
# 3. View logs:
#    sudo journalctl -u shadow-backup.service
#
# 4. Automated schedule:
#    See shadow-backup.timer (runs daily at 2 AM)
#
# ============================================================================
# Boris: "Service = what to do. Timer = when to do it.
#         Separation of concerns. SystemD architecture."
# ============================================================================

