# season-3-system-administration/episode-10-processes-systemd/solution/systemd/shadow-backup.timer
[Unit]
# KERNEL SHADOWS: Episode 10 — Backup Timer (Scheduler)
# Saint Petersburg, Russia
#
# Description: Triggers backup service daily at 2 AM
Description=KERNEL SHADOWS Backup Timer (Daily at 2 AM)
Documentation=https://github.com/kernel-shadows/episode-10
# Timer requires time-sync.target (system clock accurate)
After=time-sync.target
Requires=time-sync.target

[Timer]
# ============================================================================
# Timer Configuration
# ============================================================================

# Schedule: Daily at 2:00 AM (system local time)
OnCalendar=daily
# daily = 00:00 (midnight)
# Let's use 02:00 for less system load:
OnCalendar=*-*-* 02:00:00

# Alternative schedules:
# OnCalendar=hourly              # Every hour at :00
# OnCalendar=weekly              # Every Monday at 00:00
# OnCalendar=Mon-Fri *-*-* 02:00 # Weekdays only
# OnCalendar=*-*-* 02,14:00      # Twice daily (2 AM, 2 PM)
# OnCalendar=*-*-01 02:00        # First day of month

# Randomize start time by up to 15 minutes (prevent simultaneous load on multiple servers)
RandomizedDelaySec=15min

# If system was off when timer should have triggered, run backup on next boot
Persistent=true

# ============================================================================
# Service to Trigger
# ============================================================================
# Trigger shadow-backup.service when timer fires
Unit=shadow-backup.service

[Install]
# Enable timer to start on boot
WantedBy=timers.target

# ============================================================================
# Usage:
#
# 1. Deploy timer:
#    sudo cp shadow-backup.timer /etc/systemd/system/
#
# 2. Reload systemd:
#    sudo systemctl daemon-reload
#
# 3. Enable timer (start on boot):
#    sudo systemctl enable shadow-backup.timer
#
# 4. Start timer now:
#    sudo systemctl start shadow-backup.timer
#
# 5. Check timer status:
#    sudo systemctl status shadow-backup.timer
#
# 6. List all active timers:
#    sudo systemctl list-timers --all
#
# 7. See next run time:
#    sudo systemctl list-timers shadow-backup.timer
#
# 8. Manually trigger backup (test):
#    sudo systemctl start shadow-backup.service
#
# 9. View backup logs:
#    sudo journalctl -u shadow-backup.service -f
#
# 10. Stop timer:
#     sudo systemctl stop shadow-backup.timer
#
# 11. Disable timer (prevent start on boot):
#     sudo systemctl disable shadow-backup.timer
#
# ============================================================================
# Timer vs Cron:
#
# Cron:
#   ✓ Simple syntax: 0 2 * * * /path/to/script.sh
#   ✗ No dependency management
#   ✗ No logging integration (systemd journal)
#   ✗ No resource limits
#   ✗ If system off, job missed
#
# SystemD Timer:
#   ✓ Dependency management (After=, Requires=)
#   ✓ Integrated logging (journalctl)
#   ✓ Resource limits (CPUQuota, MemoryLimit)
#   ✓ Persistent=true catches missed runs
#   ✓ RandomizedDelaySec (load distribution)
#   ✗ More complex syntax
#
# Boris: "For modern systems, timers > cron. Better integration, more features."
#
# ============================================================================
# OnCalendar Syntax Examples:
#
# Minutely:
#   OnCalendar=*-*-* *:*:00         # Every minute
#
# Hourly:
#   OnCalendar=hourly               # Every hour at :00
#   OnCalendar=*-*-* *:15:00        # Every hour at :15
#
# Daily:
#   OnCalendar=daily                # Every day at 00:00
#   OnCalendar=*-*-* 02:30:00       # Every day at 02:30
#
# Weekly:
#   OnCalendar=weekly               # Every Monday at 00:00
#   OnCalendar=Mon *-*-* 02:00      # Every Monday at 02:00
#   OnCalendar=Fri *-*-* 18:00      # Every Friday at 18:00
#
# Monthly:
#   OnCalendar=monthly              # First of month at 00:00
#   OnCalendar=*-*-01 02:00         # First of month at 02:00
#   OnCalendar=*-*-15 02:00         # 15th of month at 02:00
#
# Multiple schedules:
#   OnCalendar=Mon-Fri *-*-* 09:00  # Weekdays at 9 AM
#   OnCalendar=Sat,Sun *-*-* 12:00  # Weekends at 12 PM
#
# ============================================================================
# Boris Kuznetsov: "Timer + Service = cron done right.
#                   Dependencies, logging, resource limits. SystemD power."
# ============================================================================

