# season-3-system-administration/episode-10-processes-systemd/solution/systemd/shadow-monitor.service
[Unit]
# KERNEL SHADOWS: Episode 10 — System Monitoring Service
# Saint Petersburg, Russia — Boris Kuznetsov's SystemD Architecture
#
# Description: Monitors Shadow operations infrastructure for suspicious activity
Description=KERNEL SHADOWS System Monitor
Documentation=https://github.com/kernel-shadows/episode-10
After=network-online.target
Wants=network-online.target
# Require logging service (journald) to be ready
Requires=systemd-journald.service
After=systemd-journald.service

[Service]
# ============================================================================
# Service Configuration
# ============================================================================
Type=simple
# oneshot: Run once and exit (for monitoring snapshots)
# simple: Run continuously (for real-time monitoring)

# User/Group to run as (NOT root unless necessary)
User=shadow-monitor
Group=shadow-monitor
# If user doesn't exist: sudo useradd -r -s /bin/false shadow-monitor

# Working Directory
WorkingDirectory=/var/operations/shadow-monitor

# ============================================================================
# Main Command
# ============================================================================
# Script that monitors system for suspicious processes, connections, files
ExecStart=/usr/local/bin/shadow-monitor.sh

# Pre-check (run before main command)
ExecStartPre=/usr/bin/test -f /usr/local/bin/shadow-monitor.sh
ExecStartPre=/usr/bin/test -x /usr/local/bin/shadow-monitor.sh

# Post-operation (cleanup)
ExecStopPost=/usr/bin/logger -t shadow-monitor "Monitoring service stopped"

# ============================================================================
# Restart Policy
# ============================================================================
# Restart automatically if crash or unexpected exit
Restart=on-failure
# Wait 10 seconds before restart
RestartSec=10s
# Max 5 restart attempts in 60 seconds
StartLimitBurst=5
StartLimitIntervalSec=60s

# ============================================================================
# Resource Limits (prevent resource exhaustion)
# ============================================================================
# CPU: Max 25% of one core (prevents DoS if service goes rogue)
CPUQuota=25%

# Memory: Max 128MB (monitoring should be lightweight)
MemoryLimit=128M
MemoryMax=128M

# Tasks: Max 50 processes/threads
TasksMax=50

# File Descriptors: Max 1024 open files
LimitNOFILE=1024

# ============================================================================
# Security Hardening (systemd security features)
# ============================================================================
# Prevent privilege escalation
NoNewPrivileges=true

# Read-only root filesystem (service cannot modify system files)
# ProtectSystem=strict
# ^ Uncomment if service only needs to write to specific locations

# Private /tmp (service cannot see other processes' temp files)
PrivateTmp=true

# Restrict network access (comment out if service needs network)
# RestrictAddressFamilies=AF_INET AF_INET6

# Prevent access to kernel modules
ProtectKernelModules=true

# Prevent access to kernel logs (dmesg)
ProtectKernelLogs=true

# Prevent access to /proc/kcore, /proc/kallsyms
ProtectKernelTunables=true

# Prevent access to /home, /root, /run/user
ProtectHome=true

# Prevent access to /sys
ProtectControlGroups=true

# Deny write to /sys, /proc/sys
ProtectProc=invisible

# Hide processes from other users
ProcSubset=pid

# ============================================================================
# Logging
# ============================================================================
# Send stdout/stderr to systemd journal
StandardOutput=journal
StandardError=journal

# Syslog identifier
SyslogIdentifier=shadow-monitor

# ============================================================================
# Environment Variables
# ============================================================================
Environment="LOG_LEVEL=INFO"
Environment="ALERT_THRESHOLD=5"
EnvironmentFile=-/etc/shadow-monitor/monitor.conf
# ^ "-" prefix = don't fail if file missing

[Install]
# Enable service to start on boot
WantedBy=multi-user.target

# ============================================================================
# Usage:
#
# 1. Deploy service unit:
#    sudo cp shadow-monitor.service /etc/systemd/system/
#
# 2. Reload systemd:
#    sudo systemctl daemon-reload
#
# 3. Enable (start on boot):
#    sudo systemctl enable shadow-monitor.service
#
# 4. Start now:
#    sudo systemctl start shadow-monitor.service
#
# 5. Check status:
#    sudo systemctl status shadow-monitor.service
#
# 6. View logs:
#    sudo journalctl -u shadow-monitor.service -f
#
# 7. Stop service:
#    sudo systemctl stop shadow-monitor.service
#
# 8. Disable (prevent start on boot):
#    sudo systemctl disable shadow-monitor.service
#
# ============================================================================
# Boris Kuznetsov: "SystemD unit = production-ready service.
#                   Resource limits, security hardening, auto-restart.
#                   That's enterprise-grade monitoring."
# ============================================================================

